# bin/sipxhomer.in
#!/bin/bash
# sipxhomer - Startup script for sipxhomer

# chkconfig: 35 85 15
# description: sipxhomer is a SIP proxy for telecommunications
# processname: sipxhomer

# Copyright (C) 2007 Pingtel Corp., certain elements licensed under a Contributor Agreement.  
# Contributors retain copyright to elements licensed under a Contributor Agreement.
# Licensed to the User under the LGPL license.

. /etc/rc.d/init.d/functions || exit 1
. @SIPX_LIBEXECDIR@/sipx-utils.sh || exit 1

prog=sipxhomer
exec="@SIPX_LIBEXECDIR@/${prog}"
pidfile="@SIPX_RUNDIR@/sipxhomer.pid"
SIPXHOMER_CONFIG="@SIPX_CONFDIR@/sipxhomer.ini"
SIPXHOMER_OPTS="--log-file @SIPX_LOGDIR@/sipxhomer.log --config-file $SIPXHOMER_CONFIG --pid-file $pidfile"

[ -e @SIPX_CONFDIR@/sysconfig/$prog ] && . @SIPX_CONFDIR@/sysconfig/$prog

start() {
    [ -x $exec ] || exit 5
    [ -f $config ] || exit 6
    if [ -n "${NoFork}" ] ; then
      runuser - @SIPXPBXUSER@ -c "$exec $SIPXHOMER_OPTS"
    else
      echo -n $"Starting $prog: "
      daemon --user @SIPXPBXUSER@ --pidfile $pidfile $exec $SIPXHOMER_OPTS
      Status=$?
      echo_success
    fi
}

stop() {
  echo -n $"Stopping sipxhomer: "
  # can take a while to shutdown because of relay
  # dont wait too long as service availability is in jeopardy
  killproc -p "${pidfile}" -d 10 "${exec}"
  Status=$?
  echo
  [ $Status -eq 0 ] && rm -f "${pidfile}"
}

restart() {
  stop
  start
}

backup() {
  for db in homer_db homer_users; do 
    mysqldump -u root ${db} > @SIPX_TMPDIR@/${db}.sql
  done
  tar -C @SIPX_TMPDIR@ -cvzf $1 homer_db.sql homer_users.sql
  rm @SIPX_TMPDIR@/homer_db.sql @SIPX_TMPDIR@/homer_users.sql
}

restore() {
  test -f ${pidfile} && running=yes
  if [ "$running" = "yes" ] ; then
    stop
  fi
  
  for db in homer_db homer_users; do 
    mysql -u root -e "drop database if exists ${db}; create database ${db};"
    tar -xzf $1 -O ${db}.sql | mysql -u root ${db}
  done

  if [ "$running" = "yes" ] ; then
    start
  fi  
}

case "$1" in
  nofork)
    NoFork=1
    start
    ;;
  start)
    start
    ;;
      
  stop)
    stop
    ;;
  restart|reload|force-reload)
    restart
    ;;
  condrestart)
    [ -f ${pidfile} ] && restart || :
    ;;
  status)
    status -p "${pidfile}" "${prog}"
    Status=$?
    ;;
  configtest)
    configtest
    ;;
  backup)
    shift
    backup $1
    ;;
  restore)
    shift
    restore $1
    ;;    
  *)
    echo "Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart|nofork|backup filename|restore filename}"
    Status=1
esac

exit $Status
